# files/system/usr/lib/systemd/system/bamos-flatpak-monitor.path
[Unit]
Description=BamOS Flatpak Monitor — triggers on flatpak install/remove
[Path]
PathModified=/var/lib/flatpak/app
[Install]
WantedBy=paths.target
