# files/system/usr/lib/systemd/system/bamos-flatpak-watcher.service
[Unit]
Description=BamOS Flatpak Watcher — monitor flatpak events for overlay integration
After=bamos-overlay-sync.service

[Service]
Type=simple
ExecStart=/usr/libexec/bamos/flatpak-event-watcher
Restart=on-failure

[Install]
WantedBy=multi-user.target
