# files/system/usr/lib/systemd/system/bamos-overlay-services.service
[Unit]
Description=BamOS Overlay Services — reenable services from overlay
After=bamos-overlay-sync.service
Before=systemd-user-sessions.service
ConditionPathExists=!/run/ostree-booted

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/bamos/bamos-overlay-services
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
