# files/system/usr/lib/systemd/system/bamos-overlay-mount.service
[Unit]
Description=BamOS Overlay Mount — mount persistent overlay on /usr
DefaultDependencies=no
After=systemd-remount-fs.service
Before=local-fs.target
ConditionPathExists=!/run/ostree-booted

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/bamos/bamos-overlay-mount
StandardOutput=journal+console

[Install]
WantedBy=local-fs.target
