# files/system/usr/lib/systemd/system/bamos-overlay-sync.service
[Unit]
Description=BamOS Overlay Sync — sync packages list into overlay
After=bamos-overlay-mount.service
Before=systemd-user-sessions.service
ConditionPathExists=!/run/ostree-booted

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/bamos/bamos-overlay-sync
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
