# files/system/usr/lib/systemd/system/bamos-cache-clean.timer
[Unit]
Description=BamOS Cache Clean — weekly dnf/podman cache cleanup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
