# files/kde-system/etc/xdg/autostart/bamos-welcome.desktop
[Desktop Entry]
Type=Application
Name=BamOS Welcome
Comment=Welcome screen for BamOS Linux — shown on first login
Exec=/usr/share/bamos/welcome/bamos-welcome.sh
Icon=bamos-logo
Terminal=false
Categories=System;Settings;
X-KDE-autostart-phase=2
X-KDE-autostart-after=panel
X-KDE-StartupNotify=false
OnlyShowIn=KDE;
X-KDE-WaitForSpeeches=false
# Show only once — disable after first successful run
X-KDE-autostart-condition=false
